# cmp-daemon.conf
start on runlevel [2345]
stop on runlevel [16]

exec /usr/sbin/cmp-daemon

respawn
respawn limit 10 5
